# MADPANDA Dark Zen/Configs/.config/hyde/themes/MADPANDA Dark Zen/helpers/systemd/user/mad-theme-watcher.service
[Unit]
Description=MADPANDA HyDE theme runtime watcher
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/mad-theme-watcher
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
